# scripts/com.artha.vault-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.artha.vault-watchdog.plist
  Artha crash-recovery watchdog — auto-encrypts sensitive state if any AI CLI session dies.

  Logic (v2 — universal AI CLI detection):
    1. Runs every 5 minutes via StartInterval
    2. Checks for .artha-decrypted lock file (signal that vault.py decrypt was run)
    3. Checks if ANY Artha-associated AI CLI is currently running:
         - Claude Code:  claude process referencing Artha
         - Gemini CLI:   gemini/node process referencing Artha
         - VS Code:      Code process with Artha folder open
    4. If lock file exists AND no active AI session → vault.py encrypt
    5. Additional TTL check: if lock is > vault_ttl_minutes old (from user_profile.yaml,
       default 30 min), auto-encrypt regardless of process state.

  Why narrowed process detection (specs/enhance.md §10.0.3):
    - Claude: `claude.*Artha` avoids false positives from unrelated claude processes
    - Gemini: `gemini.*Artha` avoids false positives from other node processes
    - VS Code: `Code.*Artha` matches VS Code with Artha folder open

  Installation (substitute your Artha checkout path before copying):
    export ARTHA_ROOT="$HOME/path/to/artha"   # e.g. ~/Library/CloudStorage/OneDrive-Personal/Artha
    sed "s|{{ARTHA_ROOT}}|${ARTHA_ROOT}|g" scripts/com.artha.vault-watchdog.plist \
      > ~/Library/LaunchAgents/com.artha.vault-watchdog.plist
    launchctl load ~/Library/LaunchAgents/com.artha.vault-watchdog.plist

  Verification:
    launchctl list | grep artha
    cat /tmp/artha-watchdog.log     (watchdog activity log)

  Ref: specs/enhance.md §10.0.3, TS §8.5
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.artha.vault-watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>
            LOCK_FILE="{{ARTHA_ROOT}}/.artha-decrypted"
            VAULT_PY="{{ARTHA_ROOT}}/scripts/vault.py"
            PYTHON="${HOME}/.artha-venvs/.venv/bin/python"
            WATCHDOG_LOG="/tmp/artha-watchdog.log"
            TIMESTAMP=$(date -Iseconds)

            # Only act if lock file exists (session was started)
            if [[ -f "${LOCK_FILE}" ]]; then

                ARTHA_SESSION_ACTIVE=false

                # Check TTL: if lock is older than 30 min, force-encrypt regardless
                LOCK_AGE=$(( $(date +%s) - $(stat -f %m "${LOCK_FILE}" 2>/dev/null || echo 0) ))
                VAULT_TTL=1800  # 30 min default (matches foundation.py LOCK_TTL)
                if [[ "${LOCK_AGE}" -gt "${VAULT_TTL}" ]]; then
                    echo "[${TIMESTAMP}] WATCHDOG: Lock TTL exceeded (age: ${LOCK_AGE}s). Auto-encrypting." >> "${WATCHDOG_LOG}"
                    if [[ -x "${PYTHON}" ]] &amp;&amp; [[ -f "${VAULT_PY}" ]]; then
                        "${PYTHON}" "${VAULT_PY}" encrypt >> "${WATCHDOG_LOG}" 2>&amp;1
                        echo "[${TIMESTAMP}] WATCHDOG: TTL auto-encrypt complete." >> "${WATCHDOG_LOG}"
                    fi
                    exit 0
                fi

                # Claude Code — narrow match to avoid false positives
                pgrep -f "claude.*Artha" &gt; /dev/null 2>&amp;1 &amp;&amp; ARTHA_SESSION_ACTIVE=true

                # Gemini CLI — narrow: gemini process working in Artha dir
                pgrep -f "gemini.*Artha" &gt; /dev/null 2>&amp;1 &amp;&amp; ARTHA_SESSION_ACTIVE=true

                # VS Code (Copilot) — VS Code with Artha folder open
                pgrep -f "Code.*Artha" &gt; /dev/null 2>&amp;1 &amp;&amp; ARTHA_SESSION_ACTIVE=true

                if [[ "${ARTHA_SESSION_ACTIVE}" == "false" ]]; then
                    # No active AI session detected — encrypt stale plaintext
                    echo "[${TIMESTAMP}] WATCHDOG: Lock file found, no active AI session. Auto-encrypting." >> "${WATCHDOG_LOG}"
                    if [[ -x "${PYTHON}" ]] &amp;&amp; [[ -f "${VAULT_PY}" ]]; then
                        "${PYTHON}" "${VAULT_PY}" encrypt >> "${WATCHDOG_LOG}" 2>&amp;1
                        echo "[${TIMESTAMP}] WATCHDOG: Auto-encrypt complete." >> "${WATCHDOG_LOG}"
                    else
                        echo "[${TIMESTAMP}] WATCHDOG: ERROR — vault.py or python not found." >> "${WATCHDOG_LOG}"
                    fi
                else
                    echo "[${TIMESTAMP}] WATCHDOG: Active AI session detected. Vault stays open." >> "${WATCHDOG_LOG}"
                fi
            fi
        </string>
    </array>

    <!-- Run every 5 minutes -->
    <key>StartInterval</key>
    <integer>300</integer>

    <!-- Start on login -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/artha-watchdog-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/artha-watchdog-stderr.log</string>

    <!-- Keep alive if it crashes -->
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
